# packaging/src/main/resources/crontab
* * * * * monitor --displayOnce | mail -s "From cron" crashubdemo@gmail.com